# scripts/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0

[program:api-bridge]
command=/opt/api-bridge/.venv/bin/python /opt/api-bridge/proxy.py
environment=PATH="/opt/api-bridge/.venv/bin:%(ENV_PATH)s",PROXY_HOST="0.0.0.0",CUSTOM_API_KEY="%(ENV_CUSTOM_API_KEY)s",CUSTOM_BASE_URL="%(ENV_CUSTOM_BASE_URL)s",DEFAULT_BACKEND="%(ENV_DEFAULT_BACKEND)s",FORCE_MODEL="%(ENV_FORCE_MODEL)s",CUSTOM_UPSTREAM_MODE="%(ENV_CUSTOM_UPSTREAM_MODE)s",INLINE_IMAGE_POLICY="%(ENV_INLINE_IMAGE_POLICY)s"
directory=/opt/api-bridge
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:claude-science]
command=/usr/local/bin/claude-science serve --no-browser --host 0.0.0.0 --port 9981
environment=ANTHROPIC_BASE_URL="http://127.0.0.1:9876"
autostart=true
autorestart=false
startretries=0
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
