# scripts/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0

[program:api-bridge]
command=/opt/api-bridge/.venv/bin/python /opt/api-bridge/proxy.py
directory=/opt/api-bridge
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:claude-science]
command=/usr/local/bin/claude-science serve --no-browser --host 0.0.0.0 --port 9981
environment=ANTHROPIC_BASE_URL="http://127.0.0.1:9876"
autostart=true
autorestart=false
startretries=0
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
